# deploy/systemd/stock-forecast-uvicorn.service
[Unit]
Description=Stock Forecast Flask App (Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=capybara27
Group=capybara27
WorkingDirectory=/home/capybara27/Documents/stock_forecast
Environment=FLASK_ENV=production
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/capybara27/Documents/stock_forecast/.venv/bin/uvicorn flask_app.asgi:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
